# services/voicebook-tts/ops/voicebook-fw.service
[Unit]
Description=voicebook-tts source allowlist (only nyla.mey.house may reach :5055)
After=docker.service
Requires=docker.service
# PartOf: a Docker daemon restart FLUSHES the DOCKER-USER chain, silently
# removing this allowlist and re-exposing the synthesis API to the whole VLAN.
# enabled-at-boot does not cover that case. PartOf ties this unit's lifecycle
# to docker.service so it is restarted with it.
PartOf=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Idempotent: delete first (ignore failure), then insert.
ExecStart=/bin/sh -c "iptables -D DOCKER-USER -p tcp --dport 5055 ! -s 10.0.20.20 -j DROP 2>/dev/null; iptables -D DOCKER-USER -p tcp --dport 5055 -s 10.0.20.20 -j RETURN 2>/dev/null; iptables -I DOCKER-USER 1 -p tcp --dport 5055 ! -s 10.0.20.20 -j DROP; iptables -I DOCKER-USER 1 -p tcp --dport 5055 -s 10.0.20.20 -j RETURN"
ExecStop=/bin/sh -c "iptables -D DOCKER-USER -p tcp --dport 5055 ! -s 10.0.20.20 -j DROP 2>/dev/null; iptables -D DOCKER-USER -p tcp --dport 5055 -s 10.0.20.20 -j RETURN 2>/dev/null; true"

[Install]
WantedBy=multi-user.target
