# config/launchd/ai.openclaw.livekit-agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Template for ai.openclaw.livekit-agent-<name>.plist.

  Substitutions (handled by scripts/deploy-agents.sh):
    {{AGENT_NAME}}            nyla | aoi | party
    {{AGENT_LABEL}}           human-readable description
    {{LANGSMITH_PROJECT}}      per-agent LangSmith project name
    {{MONOREPO_ROOT}}         absolute path to the openclaw-livekit checkout
    {{LIVEKIT_VOICE_LOGS}}    directory for voice-stack logs (default <repo>/logs/voice)

  Secrets (LIVEKIT_*, GOOGLE_API_KEY, GATEWAY_AUTH_TOKEN, DISCORD_BOT_TOKEN)
  are pulled from <repo>/secrets/livekit-agents.env by the deploy script
  (override via OPENCLAW_SECRETS env var). They are NOT committed — the
  secrets file is gitignored. Edit it and re-run scripts/deploy-agents.sh
  to rotate.
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.livekit-agent-{{AGENT_NAME}}</string>

    <key>Comment</key>
    <string>OpenClaw LiveKit agent — {{AGENT_LABEL}}</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>ProgramArguments</key>
    <array>
      <string>{{MONOREPO_ROOT}}/.venv/bin/python</string>
      <string>{{MONOREPO_ROOT}}/agents/{{AGENT_NAME}}/src/agent.py</string>
      <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{MONOREPO_ROOT}}/agents/{{AGENT_NAME}}</string>

    <key>StandardOutPath</key>
    <string>{{LIVEKIT_VOICE_LOGS}}/agent-{{AGENT_NAME}}.log</string>
    <key>StandardErrorPath</key>
    <string>{{LIVEKIT_VOICE_LOGS}}/agent-{{AGENT_NAME}}.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>{{HOME}}</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/opt/homebrew/opt/node/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>LIVEKIT_VOICE_LOGS</key>
      <string>{{LIVEKIT_VOICE_LOGS}}</string>
      <key>OPENCLAW_AGENT_NAME</key>
      <string>{{AGENT_NAME}}</string>
      <!-- OPENCLAW_BIN: rendered from env at deploy time. Falls back to
           /opt/homebrew/bin/openclaw (ARM Homebrew default). Intel Macs
           or custom install paths should set OPENCLAW_BIN before deploy. -->
      <key>OPENCLAW_BIN</key>
      <string>{{OPENCLAW_BIN}}</string>
      <key>LIVEKIT_URL</key>
      <string>{{LIVEKIT_URL}}</string>
      <key>LIVEKIT_API_KEY</key>
      <string>{{LIVEKIT_API_KEY}}</string>
      <key>LIVEKIT_API_SECRET</key>
      <string>{{LIVEKIT_API_SECRET}}</string>
      <key>GOOGLE_API_KEY</key>
      <string>{{GOOGLE_API_KEY}}</string>
      <key>GATEWAY_AUTH_TOKEN</key>
      <string>{{GATEWAY_AUTH_TOKEN}}</string>
      <key>DISCORD_BOT_TOKEN</key>
      <string>{{DISCORD_BOT_TOKEN}}</string>
      <!-- Canonical Musubi (v0.6.0+) — base URL shared by every agent,
           per-agent bearer token resolved by scripts/deploy-agents.sh
           from secrets/livekit-agents.env (MUSUBI_V2_TOKEN_NYLA /
           _AOI). Party reuses Nyla's token — see PartyAgent config. -->
      <key>MUSUBI_V2_BASE_URL</key>
      <string>{{MUSUBI_V2_BASE_URL}}</string>
      <key>MUSUBI_V2_TOKEN</key>
      <string>{{MUSUBI_V2_TOKEN}}</string>
      <!-- Chained-pipeline keys (Party only — Whisper STT + ElevenLabs TTS).
           Other agents won't use these but it costs nothing to have them in
           the env; rendering them here keeps the template uniform across
           agent plists. -->
      <key>OPENAI_API_KEY</key>
      <string>{{OPENAI_API_KEY}}</string>
      <key>ELEVENLABS_API_KEY</key>
      <string>{{ELEVENLABS_API_KEY}}</string>
      <!-- OpenTelemetry tracing (vendor-neutral). Master switch:
           OPENCLAW_OTEL_ENABLED=true (or legacy LANGSMITH_TRACING=true).
           Every LiveKit framework span and outbound HTTP call ships to
           every backend in OPENCLAW_OTEL_EXPORTERS. -->
      <key>OPENCLAW_OTEL_ENABLED</key>
      <string>{{OPENCLAW_OTEL_ENABLED}}</string>
      <key>OPENCLAW_OTEL_EXPORTERS</key>
      <string>{{OPENCLAW_OTEL_EXPORTERS}}</string>
      <key>OPENCLAW_OTEL_DEBUG</key>
      <string>{{OPENCLAW_OTEL_DEBUG}}</string>
      <key>OPENCLAW_OTEL_HTTP_INSTRUMENTATION</key>
      <string>{{OPENCLAW_OTEL_HTTP_INSTRUMENTATION}}</string>
      <key>OPENCLAW_DEPLOYMENT_ENVIRONMENT</key>
      <string>{{OPENCLAW_DEPLOYMENT_ENVIRONMENT}}</string>
      <key>OPENCLAW_SERVICE_VERSION</key>
      <string>{{OPENCLAW_SERVICE_VERSION}}</string>
      <key>LANGSMITH_TRACING</key>
      <string>{{LANGSMITH_TRACING}}</string>
      <key>LANGSMITH_PROJECT</key>
      <string>{{LANGSMITH_PROJECT}}</string>
      <key>LANGSMITH_ATTACH_AUDIO</key>
      <string>{{LANGSMITH_ATTACH_AUDIO}}</string>
      <key>LIVEKIT_EGRESS_HOST_RECORDINGS_DIR</key>
      <string>{{LIVEKIT_EGRESS_HOST_RECORDINGS_DIR}}</string>
      <key>LIVEKIT_EGRESS_CONTAINER_RECORDINGS_DIR</key>
      <string>{{LIVEKIT_EGRESS_CONTAINER_RECORDINGS_DIR}}</string>
      <key>OTEL_EXPORTER_OTLP_ENDPOINT</key>
      <string>{{OTEL_EXPORTER_OTLP_ENDPOINT}}</string>
      <key>OTEL_EXPORTER_OTLP_HEADERS</key>
      <string>{{OTEL_EXPORTER_OTLP_HEADERS}}</string>
      <key>OPENCLAW_PHOENIX_ENDPOINT</key>
      <string>{{OPENCLAW_PHOENIX_ENDPOINT}}</string>
      <key>OPENCLAW_PHOENIX_HEADERS</key>
      <string>{{OPENCLAW_PHOENIX_HEADERS}}</string>
      <key>OPENCLAW_OTLP_ENDPOINT</key>
      <string>{{OPENCLAW_OTLP_ENDPOINT}}</string>
      <key>OPENCLAW_OTLP_HEADERS</key>
      <string>{{OPENCLAW_OTLP_HEADERS}}</string>
      <key>OPENCLAW_OTEL_LOGS_ENABLED</key>
      <string>{{OPENCLAW_OTEL_LOGS_ENABLED}}</string>
      <key>OPENCLAW_OTLP_LOGS_ENDPOINT</key>
      <string>{{OPENCLAW_OTLP_LOGS_ENDPOINT}}</string>
      <key>OPENCLAW_OTLP_LOGS_HEADERS</key>
      <string>{{OPENCLAW_OTLP_LOGS_HEADERS}}</string>
      <key>LANGSMITH_VERBOSE_TELEMETRY</key>
      <string>{{LANGSMITH_VERBOSE_TELEMETRY}}</string>
      <key>LANGSMITH_PROCESSOR_DEBUG</key>
      <string>{{LANGSMITH_PROCESSOR_DEBUG}}</string>
      <key>NODE_EXTRA_CA_CERTS</key>
      <string>/etc/ssl/cert.pem</string>
    </dict>
  </dict>
</plist>
